# containers/cbackup/container/etc/supervisor/conf.d/cbackup.conf
[program:cbackup]
command=/opt/cbackup/bin/cbackup.jar
autorestart = true
stdout_logfile = /var/log/cbackup/cbackup-supervisor.log
redirect_stderr = true