# containers/cbackup/cbackup/container/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autorestart = true
stdout_logfile = /var/log/sshd/sshd-supervisor.log
redirect_stderr = true