# containers/bginetaccess/container/etc/supervisor/conf.d/bginetaccess.conf
[program:bginetaccess]
command=/opt/bgbilling/BGInetAccess/bin/application.sh start
autorestart = true
stdout_logfile = /var/log/bginetaccess/bginetaccess-supervisor.log
redirect_stderr = true