# containers/bgbilling/bginetaccess/default/8.2110-2111301445/BGInetAccess/script/bginet_access.service
[Unit]
Description=BGInetAccess
After=network.target activemq.service
Requires=activemq.service

[Service]
Type=forking
TimeoutStartSec=120s
ExecStart=/opt/bgbilling/BGInetAccess/bin/application.sh start
ExecStop=/opt/bgbilling/BGInetAccess/bin/application.sh stop

[Install]
WantedBy=multi-user.target
